# linuxhub.service
[Unit]
Description=LinuxHub Application
Documentation=https://github.com/oguzkaganeren/linuxhub

[Install]
WantedBy=default.target
